# conf/server-monitor.timer
[Unit]
Description=Execute service based on timer
[Timer]
OnBootSec=900
OnUnitActiveSec=900
AccuracySec=1ms
[Install]
WantedBy=timers.target